# redissyncer.service
[Service]
ExecStart=/home/redissyncer/redissyncer/startcompose.sh

[Install]
WantedBy=default.target
